# server/bin/blackhole.init.plist
<?xml version="1.0" encoding="utf-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
      <key>Label</key>
      <string>blackhole.init</string>
      <key>ProgramArguments</key>
      <array>
        <string>/usr/local/blackhole/blackhole-start.sh</string>
      </array>
	<key>KeepAlive</key>
	<false/>
      <key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/blackhole.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/blackhole.out</string>
	<key>UserName</key>
	<string>root</string>
  </dict>
</plist>

